# com.quietspaceclub.blog-autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Label must be unique across all launchd agents -->
  <key>Label</key>
  <string>com.quietspaceclub.blog-autopush</string>

  <!-- Run autopush.sh via bash -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/petehallett/Documents/websites/quiet-space-neuro-index/quiet-space-neuro-index/autopush.sh</string>
  </array>

  <!-- Poll every 300 seconds (5 minutes) -->
  <key>StartInterval</key>
  <integer>300</integer>

  <!-- Don't run immediately on load — wait for first interval -->
  <key>RunAtLoad</key>
  <false/>

  <!-- Log output to repo root so it's easy to inspect -->
  <key>StandardOutPath</key>
  <string>/Users/petehallett/Documents/websites/quiet-space-neuro-index/quiet-space-neuro-index/autopush.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/petehallett/Documents/websites/quiet-space-neuro-index/quiet-space-neuro-index/autopush.log</string>

  <!-- Keep environment variables needed for git -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/petehallett</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
